# etc/systemd/system/puavo-vpn-internal.service
[Unit]
Description=Simple WSGI Server

[Service]
Type=simple
ExecStart=/usr/local/sbin/puavo-vpn-internal

[Install]
WantedBy=multi-user.target